# scripts/com.lei.camwatch-logrotate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lei.camwatch-logrotate</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/lei/github/camwatch/scripts/camwatch_log_rotate.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>17</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/lei/Library/Logs/camwatch/logrotate.out</string>

    <key>StandardErrorPath</key>
    <string>/Users/lei/Library/Logs/camwatch/logrotate.err</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
